# backup-auto.service
[Unit]
Description=Backup automatico via rclone

[Service]
ExecStart=%h/.local/scripts/backup-auto.sh


