# backup-auto.service
[Unit]
Description=Automatic backup via rclone

[Service]
ExecStart=%h/.local/scripts/backup-auto.sh


