# cleanup-versions.service
[Unit]
Description=Cleaning old versions on remote backup

[Service]
ExecStart=%h/.local/scripts/cleanup-versions.sh


